# infrastructure/ansible/supervisor/templates/disque.conf
[program:disque]
command=disque-server
process_name=%(program_name)s

directory=/home/ansible

stopsignal=TERM

autostart=true
autorestart=true
