# infrastructure/ansible/supervisor/templates/server.conf
[program:server]
command=/usr/bin/python3 server.py
process_name=%(program_name)s

directory=/home/ansible/app

stopsignal=TERM

autostart=true
autorestart=true

environment=PRODUCTION="Yes"
