# infrastructure/ansible/supervisor/templates/server.conf
[program:server]
command=python3 manage.py runserver
process_name=%(program_name)s

directory=/home/ansible

stopsignal=TERM

autostart=true
autorestart=true

environment=PRODUCTION="Yes"

