# infrastructure/ansible/supervisor/templates/server.conf
[program:server]
command=gunicorn linkalytics.wsgi:app -c /home/ansible/linkalytics/gunicorn.conf.py
process_name=%(program_name)s

directory=/home/ansible

stopsignal=TERM

autostart=true
autorestart=true

environment=PRODUCTION="Yes"

