# infrastructure/ansible/supervisor/templates/tika.conf
[program:tika]
command=java -jar tika-server.jar -h 0.0.0.0 -p 9998
process_name=%(program_name)s

directory=/home/ansible

stopsignal=TERM

autostart=true
autorestart=true
